# host/upstart.conf
description "Drycc host daemon"

start on runlevel [2345]
stop on runlevel [!2345]
respawn
respawn limit 100 60

exec /usr/local/bin/drycc-host daemon
